# docker/php/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770

[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[program:queue-worker-redis]
process_name=redis_worker_%(process_num)02d
command=php /var/www/html/artisan queue:work redis --sleep=3 --tries=3 --timeout=600 --max-time=3600
autostart=true
autorestart=true
numprocs=3
redirect_stderr=true
stdout_logfile=/var/log/supervisor/redis_worker_%(process_num)02d.log
stdout_logfile_maxbytes=5MB
stopwaitsecs=600

; ── Database workers: drain jobs dispatched before QUEUE_CONNECTION=redis ──
; Once jobs table is empty, these will idle harmlessly.
; Safe to remove this block after the backlog is fully cleared.
[program:queue-worker-database]
process_name=db_worker_%(process_num)02d
command=php /var/www/html/artisan queue:work database --sleep=3 --tries=3 --timeout=600 --max-time=3600
autostart=true
autorestart=true
numprocs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/db_worker_%(process_num)02d.log
stdout_logfile_maxbytes=5MB
stopwaitsecs=600
